# mes_orb_bot/deploy/systemd/mes-orb-dashboard.service
[Unit]
Description=MES ORB Bot Dashboard (Streamlit)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__BOT_USER__
WorkingDirectory=__BOT_ROOT__
Environment=PYTHONUNBUFFERED=1

ExecStart=__BOT_ROOT__/.venv/bin/streamlit run __BOT_ROOT__/dashboard.py \
    --server.port 8502 \
    --server.address 127.0.0.1 \
    --server.headless true \
    --server.fileWatcherType none \
    --browser.gatherUsageStats false

Restart=always
RestartSec=10
TimeoutStopSec=15

Nice=15
SyslogIdentifier=mes-orb-dashboard

[Install]
WantedBy=multi-user.target
